# server/supervisor/development.ini
[program:varify-development]
directory = /home/devuser/webapps/varify-env/varify
command = /home/devuser/webapps/varify-env/bin/uwsgi --ini /home/devuser/webapps/varify-env/varify/server/uwsgi/development.ini
autostart = true
autorestart = true
