# server/supervisor/local.ini
[program:varify-local]
directory = /home/devuser/webapps/varify-env/varify
command = /home/devuser/webapps/varify-env/bin/uwsgi --ini /home/devuser/webapps/varify-env/varify/server/uwsgi/local.ini
autostart = true
autorestart = true
